# crontab
0 2 * * * backup_date=$(date +\%Y\%m\%d)
mkdir -p /backup/$backup_date
mysqldump -usantiago.etchebarne -p43375773 project > "/backup/$backup_date/mysqldump_$backup_date.sql"